# etc/fumobot.path
[Unit]
Description=Fumbot Update

[Path]
PathChanged=/usr/local/fumobot/bin/linux-x64/fumobot
Unit=fumobot-restart.service

[Install]
WantedBy=multi-user.target