# etc/fumobot-restart.service
[Unit]
Description=Fumobot Restart

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart fumobot.service